# deploy/canned-agent.service.example
[Unit]
Description=Canned Agent
After=network.target

[Service]
Type=simple
User=drandall
Group=drandall
WorkingDirectory=/home/drandall/CannedAgent
EnvironmentFile=-/home/drandall/CannedAgent/deploy/canned-agent.env
ExecStart=/home/drandall/CannedAgent/venv/bin/uvicorn src.app:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
